# elk_config/templates/kibana.service.j2
[Unit]
Description=Kibana service
After=network.target

[Service]
Type=simple
ExecStart=/data/{{ kibana }}/bin/kibana

[Install]
WantedBy=multi-user.target
